# build_deb/postinst
#!/bin/bash
set -e

# Ensure config directory and config file permissions are secure for root-run service
mkdir -p /etc/nest
if [ -f /etc/nest/nest.conf ]; then
    chown root:root /etc/nest/nest.conf
    chmod 600 /etc/nest/nest.conf
fi

# Load variables from /etc/nest/nest.conf
if [ -f /etc/nest/nest.conf ]; then
    set -a
    # shellcheck disable=SC1091
    source /etc/nest/nest.conf
    set +a
fi

# Use the PLUGINS_DIR variable from config and create the directory
mkdir -p "${PLUGINS_DIR:-/opt/nest-plugins}"

# Install systemd service
systemctl daemon-reload
systemctl enable nest.service

# Start the service
systemctl start nest.service